# etc/systemd/system/systemd-modules-load.service.d/update-deps.conf
[Service]
ExecStartPre=-/bin/sh -c 'md="/lib/modules/$(uname -r)"; test -z "$(find "$md/" -newer "$md/modules.dep" -name "*.ko")" || depmod'
