# etc/systemd/system/NetworkManager.service.d/conflict-dhcpcd.conf
[Unit]
Conflicts=dhcpcd.service
Before=dhcpcd.service

