# etc/systemd/system/systemd-journald-audit.socket.d/explicit-enable.conf
[Unit]
ConditionPathExists=/etc/audit/journald-enabled
