# etc/systemd/user/save-setting@.path
[Path]
PathChanged=%f

[Install]
WantedBy=default.target
