# etc/systemd/user/xsession.target
[Unit]
Description=Started last from /etc/X11/Xsession.d
