# etc/systemd/system/docker.service.d/need-start-file.conf
[Unit]
ConditionPathExists=/etc/default/docker.start-by-default
