# etc/systemd/system/docker.socket.d/need-docker-dir.conf
[Unit]
ConditionPathExists=/var/lib/docker
