# etc/systemd/system/ssh.socket.d/host-keys.conf
[Socket]
ExecStartPre=/bin/sh -c 'test -n "$(find /etc/ssh/ -name "ssh_host_*_key")" || dpkg-reconfigure openssh-server'
