# etc/systemd/system/dkms-autoinstall.service
[Unit]
Description=Run dkms autoinstall

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/dkms autoinstall
