# etc/systemd/system/run-live.service
[Unit]
ConditionKernelCommandLine=run-live
ConditionFileIsExecutable=/live/.run
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/live/.run
