# etc/systemd/system/simplehttpserver.service
[Service]
User=www-data
WorkingDirectory=/var/www
ExecStart=/usr/bin/python -mSimpleHTTPServer
