# etc/systemd/system/tor@.service.d/require-enabled.conf
[Unit]
Wants=check-tor-enabled.service
After=check-tor-enabled.service
ConditionPathExists=/var/lib/tor/.enabled
