# etc/systemd/system/tor@.service.d/require-enabled.conf
[Unit]
Requires=check-tor-enabled.service
After=check-tor-enabled.service
